# raspberry_pi/roles/mackerel-agent/templates/mackerel-agent.service.j2
[Unit]
Description = mackerel-agent

[Service]
ExecStart = /usr/local/mackerel-agent/mackerel-agent
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target

